# auth-server/wk-auth-server.service
[Unit]
Description=Workit OAuth Callback Server
After=network.target
Documentation=https://github.com/namastexlabs/workit

[Service]
Type=simple
User=genie
Group=genie
WorkingDirectory=/home/genie/workspace/repos/workit/auth-server

# Load credentials from env file
EnvironmentFile=/home/genie/.config/workit/credentials.env

# Server config
Environment=PORT=8089
Environment=REDIRECT_URL=https://wkauth.namastex.io/callback

ExecStart=/home/genie/workspace/repos/workit/auth-server/auth-server \
    --port=${PORT} \
    --redirect-url=${REDIRECT_URL}

Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target
